# deploy/wiki-gap.service
[Unit]
Description=wiki-gap dashboard (FastAPI + Uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tama
Group=tama
WorkingDirectory=/home/tama/wiki-gap
Environment="PYTHONPATH=/home/tama/wiki-gap"
EnvironmentFile=/home/tama/wiki-gap/.env
# Tailscale IP 限定 bind: public IP には漏らさない、Tailnet 内 (Mac/iPhone等) からは見れる
ExecStart=/home/tama/wiki-gap/.venv/bin/uvicorn src.web.app:app --host ${WIKI_GAP_BIND_HOST} --port ${WIKI_GAP_BIND_PORT}
Restart=on-failure
RestartSec=10s

# セキュリティ
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/tama/wiki-gap/data /home/tama/wiki-gap/logs

[Install]
WantedBy=multi-user.target
