# deploy/wiki-gap-crawl.timer
[Unit]
Description=Monthly wiki-gap crawl @ JST 03:00 (Wikimedia の負荷の薄い時間帯)

[Timer]
# 毎月 1 日 JST 03:00 に実行 (旧: 毎日。2026-05-27 に月1へ変更 — 実利用が月数回ペースのため)
# サーバー TZ が JST であることを前提 (UTC ベースなら 18:00 で書く)
OnCalendar=*-*-01 03:00:00
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
